# skel/lib/systemd/system/dcache.target
[Unit]
Description=dCache target to start/stop all dcache@*.service instances

[Install]
WantedBy=default.target