# packages/fhs/src/main/deb/postinst
#!/bin/sh -e

if [ "$1" = configure ]; then
    # Make sure the administrative user exists
    if ! getent passwd dcache > /dev/null; then
        adduser --system --quiet --home /var/lib/dcache --no-create-home \
            --shell /bin/bash --group --gecos "dCache administrator" dcache
        adduser --quiet dcache ssl-cert
    fi

    # check validity of dcache user and group
    if [ "`id -u dcache`" -eq 0 ]; then
        echo "The dcache system user must not have uid 0 (root).
Please fix this and reinstall this package." >&2
        exit 1
    fi
    if [ "`id -g dcache`" -eq 0 ]; then
        echo "The dcache system user must not have root as primary group.
Please fix this and reinstall this package." >&2
        exit 1
    fi

    # ensure directory ownership
    chown dcache:dcache /var/lib/dcache
    chown dcache:dcache /var/lib/dcache/alarms
    chown dcache:dcache /var/lib/dcache/config
    chown dcache:dcache /var/lib/dcache/billing
    chown dcache:dcache /var/lib/dcache/httpd
    chown dcache:dcache /var/lib/dcache/plots
    chown dcache:dcache /var/lib/dcache/pool-history
    chown dcache:dcache /var/lib/dcache/statistics
    chown dcache:dcache /var/lib/dcache/credentials
    chown dcache:dcache /var/lib/dcache/star
    chown dcache:dcache /var/log/dcache
    chown dcache:dcache /var/spool/dcache/star

    # protect directories that only dCache should access
    chmod 700 /var/lib/dcache/alarms
    chmod 700 /var/lib/dcache/credentials
    chmod 700 /var/lib/dcache/httpd
    chmod 700 /var/lib/dcache/plots
    chmod 700 /var/lib/dcache/pool-history
    chmod 700 /var/lib/dcache/resilience
    chmod 700 /var/lib/dcache/statistics

    # allow somewhat more relaxed permissions elsewhere
    chmod 750 /var/lib/dcache/billing
    chmod 770 /var/lib/dcache/star
    chmod 755 /var/spool/dcache/star

    # generate admin door ssh2 server key
    if [ ! -f /etc/dcache/admin/ssh_host_rsa_key ]; then
        ssh-keygen -q -t rsa -f /etc/dcache/admin/ssh_host_rsa_key -N ""
        chmod 640 /etc/dcache/admin/ssh_host_rsa_key
        chgrp dcache /etc/dcache/admin/ssh_host_rsa_key
    fi

    # decide how we run dcache domains
    if [ ! -f /usr/share/dcache/lib/services.sh ]; then
        if type "systemctl" > /dev/null 2>&1; then
            ln -sf /usr/share/dcache/lib/services-systemd.sh /usr/share/dcache/lib/services.sh
        else
            ln -sf /usr/share/dcache/lib/services-daemon.sh /usr/share/dcache/lib/services.sh
        fi
    fi

    if type "systemctl" > /dev/null 2>&1; then
        sed -i /usr/share/dcache/defaults/dcache.properties -e "s/(not-for-services)dcache.log.format.file=.*/(not-for-services)dcache.log.format.file=%-5level %-20([%X{cells.cell}]) %m %replace(<%X{org.dcache.ndc}>){'^<>\$',''}%n/"
        rm -f /etc/init.d/dcache
        update-rc.d dcache remove
        systemctl daemon-reload
    fi

    service rsyslog restart
fi

#DEBHELPER#
