# packages/fhs/src/main/deb/prerm
#!/bin/sh -e

if type "systemctl" > /dev/null 2>&1; then
        systemctl stop dcache@*
else
        invoke-rc.d dcache stop || :
fi

#DEBHELPER#
