# skel/bin/dcache-srm
#!/bin/sh
# $Id: dcache-srm,v 1.5 2007-08-14 21:13:00 tigran Exp $
#
# chkconfig: - 91 9
# description: srm tomcat startup script
#

ourHomeDir=/opt/d-cache

cmd1=$1

if [ -r ${ourHomeDir}/etc/srm_setup.env ] ; then
  . ${ourHomeDir}/etc/srm_setup.env
fi

hostname=`hostname | awk -F. '{print $1}'`

export JAVA_HOME
export JAVA_OPTS
export CATALINA_HOME=${TOMCAT_BASE_PATH}/apache-tomcat-5.5.20
export CATALINA_PID=${ourHomeDir}/config/lastPid.srm-${hostname}

export X509_CERT
export X509_KEY
export SRM_V2_PORT


case $cmd1 in
  start)
        ${TOMCAT_BASE_PATH}/apache-tomcat-5.5.20/bin/startup.sh
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && touch /var/lock/subsys/srm-tomcat
	echo "Pinging srm server to wake it up, will take few seconds ..."
	export CLASSPATH
	sleep 2 
	for jar in `find /opt/d-cache/classes/ -name '*jar'`; do  CLASSPATH=${CLASSPATH}:${jar}; done
	$JAVA_HOME/bin/java org.dcache.srm.client.SrmStartUpPing
	echo "Done"
        ;;
  stop)
         ${TOMCAT_BASE_PATH}/apache-tomcat-5.5.20/bin/shutdown.sh
          RETVAL=$?
          echo
          [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/srm-tomcat
          ;;
  restart|reload)
        $0 stop
        $0 start
        RETVAL=$?
        ;;
  *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
esac
exit 0

